# gpio/setup_buttons.service
[Unit]
ConditionFileIsExecutable=/bin/setup_buttons

[Service]
Type=simple
ExecStart=/bin/setup_buttons

[Install]
WantedBy=basic.target
