# gpio/watch_buttons.service
[Unit]
ConditionFileIsExecutable=/bin/watch_buttons

[Service]
Type=simple
ExecStart=/bin/watch_buttons

[Install]
WantedBy=multi-user.target
